# skel/026-cluster-consul/etc/NetworkManager/dispatcher.d/50-update-dnsmasq-resolv-conf
#!/bin/bash

## script to be executed when dhcp changes are made; maintains /etc/resolv.conf
## man NetworkManager

interface="${1}"
action="${2}"

/bin/logger -t "${0}" \
    "invoked for interface ${interface:-<not_provided>} and action ${action} with nameservers '${IP4_NAMESERVERS}', domains '${IP4_DOMAINS}'"

if [ -n "${IP4_NAMESERVERS}" ] && [ -n "${IP4_DOMAINS}" ]; then
    echo "# generated by ${0}; do not edit" >| /etc/resolv.conf.dnsmasq
    
    for ns in ${IP4_NAMESERVERS}; do
        echo "nameserver ${ns}" >> /etc/resolv.conf.dnsmasq
    done

    ## use IP of eth0, so we can use the same resolv.conf and mount it in
    ## containers
    local_ip=$( ip addr show eth0 | awk '/inet / {print substr($2, 0, index($2, "/") - 1)}' )
    echo "# generated by ${0}; do not edit" >| /etc/resolv.conf
    echo "nameserver ${local_ip:-127.0.0.1}" >> /etc/resolv.conf
    
    if [ -n "${IP4_DOMAINS}" ]; then
        echo "search ${IP4_DOMAINS}" >> /etc/resolv.conf
    fi
fi
